# initd/systemd.gnugk.service
[Unit]
Description=GNU Gatekeeper

[Service]
LimitRTPRIO=100

